# packaging/systemd/user/qcma.service.in
[Unit]
Description="Qcma daemon"
Documentation=man:qcma_cli(6)
ConditionPathExists=%h/.config/codestation/qcma.conf

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/qcma_cli --verbose
ExecReload=/bin/kill -HUP $MAINPID

[Install]
Alias=qcma_cli.service
WantedBy=default.target
